# etc/cron.d/feed2json.org
MAILTO=andychilton@gmail.com

# mins hour dofm mnth dofw   user          command
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
     0    *    *    *    *   chilts        find /var/cache/feed2json/feed/ -type f -mmin +60 -delete
# ---- ---- ---- ---- ----   -----------   ----------------------------------------------------------------------------
